# leanote_standalone/supervisord.conf
[supervisord]
nodaemon=true
[program:mongodb]
command=mongod --dbpath /db
[program:leanote]
command=/bin/bash -c "/init.sh && /leanote/bin/run.sh"
